# debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        systemctl stop travel-net.service || true
        systemctl disable travel-net.service || true
        ;;
esac
